# mohlenhoff_alpha2/rootfs/etc/services.d/alpha2-webui/finish
#!/command/with-contenv bashio
# ==============================================================================
# Take down the S6 supervision tree when Alpha 2 Integration fails
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
  bashio::log.warning "Alpha 2 Integration crashed, halting add-on"
  /run/s6/basedir/bin/halt
fi

bashio::log.info "Alpha 2 Integration stopped, restarting..."