# mohlenhoff_alpha2/rootfs/etc/services.d/alpha2-webui/run
#!/command/with-contenv bashio
# ==============================================================================
# Start Alpha 2 Web UI service
# ==============================================================================

# Set ingress port and URI as environment variables
export INGRESS_PORT=$(bashio::addon.ingress_port)
export INGRESS_URI=$(bashio::addon.ingress_entry)

# Start the web UI
bashio::log.info "Starting Alpha 2 Web UI on port ${INGRESS_PORT} with URI ${INGRESS_URI}..."
exec python3 /usr/bin/alpha2_webui.py